# services/chaotic-sort-logs.timer
[Unit]
Description=Chaotic's automatic log sorting utility

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
