# services/chaotic-tkg-wine.timer
[Unit]
Description=Chaotic+TkG's wine routine

[Timer]
OnCalendar=*-*-* 15:20:00
Persistent=false

[Install]
WantedBy=timers.target
