# services/chaotic-hourly.timer
[Unit]
Description=Chaotic's hourly routine

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
